# src/bot_api_v1/scripts/c.conf


# 通用环境变量 (将实际值填入下方每个 program 的 environment 行)
# 注意: 多个变量用逗号分隔, 变量值用引号包起来
# 请替换掉 <YOUR_STRONG_SECRET_KEY>, <YOUR_APP_LOG_PATH>, <YOUR_WHISPER_MODEL>
# 并确认 Redis/DB 密码和用户名
# PYTHONPATH="/code/bot_app/bot_api_v1/src" # 假设项目根目录是 /code/bot_app/bot_api_v1
# ENVIRONMENT="production"
# SECRET_KEY="x95cYCDGkKeKCiL1lYzi3jEUbnzmeKanNVGJkZm_zC0"
# JWT_SECRET_KEY="S0Zk0cjmbG0xf39Wct-NI1d6EOTeW_T00vcON61-IBA"
# DATABASE_URL="postgresql+asyncpg://cappa_rw:RWcappaDb!!!2025@10.0.16.12:5432/cappadocia_v1"
# CELERY_BROKER_URL="redis://:login4RDS!!!@10.0.16.12:6379/0"
# CELERY_RESULT_BACKEND="redis://:login4RDS!!!@10.0.16.12:6379/1"
# CACHE_REDIS_URL="redis://:login4RDS!!!@10.0.16.12:6379/2"
# SHARED_TEMP_DIR="/srv/nfs/shared" # NFS 源目录路径
# LOG_FILE_PATH="/var/log/bot_api_v1" # 例如 /var/log/bot_api_v1
# LOG_LEVEL="INFO"
# DEBUG="False"
# WHISPER_MODEL="base" # 例如 base

[program:fastapi_app]
# !! 修改 command 中的 uvicorn 启动目标为你实际的 app factory 或 app 实例 !!
# !! workers=2 是基于 2核4G 服务器的建议值，可根据负载调整 !!
command=/code/bot_app/bot_api_v1/venv/bin/python -m uvicorn bot_api_v1.app.core.app_factory:create_app --host 0.0.0.0 --port 8083 --workers 2 --loop uvloop --http httptools
directory=/code/bot_app/bot_api_v1 # 确认项目根目录
user=ubuntu
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/fastapi_app_stdout.log
stderr_logfile=/var/log/supervisor/fastapi_app_stderr.log
# !! 将上面注释掉的通用环境变量复制并取消注释到下面这一行, 并填入实际值 !!
environment=PYTHONPATH="/code/bot_app/bot_api_v1/src",ENVIRONMENT="production",SECRET_KEY="x95cYCDGkKeKCiL1lYzi3jEUbnzmeKanNVGJkZm_zC0",JWT_SECRET_KEY="S0Zk0cjmbG0xf39Wct-NI1d6EOTeW_T00vcON61-IBA",DATABASE_URL="postgresql+asyncpg://cappa_rw:RWcappaDb!!!2025@10.0.16.12:5432/cappadocia_v1",CELERY_BROKER_URL="redis://:login4RDS!!!@10.0.16.12:6379/0",CELERY_RESULT_BACKEND="redis://:login4RDS!!!@10.0.16.12:6379/1",CACHE_REDIS_URL="redis://:login4RDS!!!@10.0.16.12:6379/2",SHARED_TEMP_DIR="/srv/nfs/shared",LOG_FILE_PATH="/var/log/bot_api_v1",LOG_LEVEL="INFO",DEBUG="False",WHISPER_MODEL="base"


[program:celery_worker_A]
# !! 修改 command 中的 -A 参数为你实际的 Celery App 实例路径 !!
# !! -c 1 是基于 2核4G 服务器的保守建议值，处理媒体/默认任务可能也耗资源 !!
command=/code/bot_app/bot_api_v1/venv/bin/celery -A bot_api_v1.app.tasks.celery_app worker -l info -Q celery,media_extraction -c 1 -n taskA
directory=/code/bot_app/bot_api_v1 # 确认项目根目录
user=ubuntu
autostart=true
autorestart=true
stopwaitsecs = 600
stopsignal=TERM
stdout_logfile=/var/log/supervisor/celery_worker_A_stdout.log
stderr_logfile=/var/log/supervisor/celery_worker_A_stderr.log
# !! 将上面注释掉的通用环境变量复制并取消注释到下面这一行, 并填入实际值 !!
environment=PYTHONPATH="/code/bot_app/bot_api_v1/src",ENVIRONMENT="production",SECRET_KEY="x95cYCDGkKeKCiL1lYzi3jEUbnzmeKanNVGJkZm_zC0",JWT_SECRET_KEY="S0Zk0cjmbG0xf39Wct-NI1d6EOTeW_T00vcON61-IBA",DATABASE_URL="postgresql+asyncpg://cappa_rw:RWcappaDb!!!2025@10.0.16.12:5432/cappadocia_v1",CELERY_BROKER_URL="redis://:login4RDS!!!@10.0.16.12:6379/0",CELERY_RESULT_BACKEND="redis://:login4RDS!!!@10.0.16.12:6379/1",CACHE_REDIS_URL="redis://:login4RDS!!!@10.0.16.12:6379/2",SHARED_TEMP_DIR="/srv/nfs/shared",LOG_FILE_PATH="/var/log/bot_api_v1",LOG_LEVEL="INFO",DEBUG="False",WHISPER_MODEL="base"


[program:celery_worker_log]
# !! 修改 command 中的 -A 参数为你实际的 Celery App 实例路径 !!
# !! -c 1 对于日志处理通常足够 !!
command=/code/bot_app/bot_api_v1/venv/bin/celery -A bot_api_v1.app.tasks.celery_app worker -l info -Q logging -c 1 -n task_log
directory=/code/bot_app/bot_api_v1 # 确认项目根目录
user=ubuntu
autostart=true
autorestart=true
stopwaitsecs = 600
stopsignal=TERM
stdout_logfile=/var/log/supervisor/celery_worker_log_stdout.log
stderr_logfile=/var/log/supervisor/celery_worker_log_stderr.log
# !! 将上面注释掉的通用环境变量复制并取消注释到下面这一行, 并填入实际值 !!
environment=PYTHONPATH="/code/bot_app/bot_api_v1/src",ENVIRONMENT="production",SECRET_KEY="x95cYCDGkKeKCiL1lYzi3jEUbnzmeKanNVGJkZm_zC0",JWT_SECRET_KEY="S0Zk0cjmbG0xf39Wct-NI1d6EOTeW_T00vcON61-IBA",DATABASE_URL="postgresql+asyncpg://cappa_rw:RWcappaDb!!!2025@10.0.16.12:5432/cappadocia_v1",CELERY_BROKER_URL="redis://:login4RDS!!!@10.0.16.12:6379/0",CELERY_RESULT_BACKEND="redis://:login4RDS!!!@10.0.16.12:6379/1",CACHE_REDIS_URL="redis://:login4RDS!!!@10.0.16.12:6379/2",SHARED_TEMP_DIR="/srv/nfs/shared",LOG_FILE_PATH="/var/log/bot_api_v1",LOG_LEVEL="INFO",DEBUG="False",WHISPER_MODEL="base"


sudo supervisorctl reread
sudo supervisorctl update
sudo supervisorctl start all # 尝试启动所有新配置的程序


sudo supervisorctl status all

sudo tail -n 50 /var/log/supervisor/fastapi_app_stderr.log
sudo tail -n 50 /var/log/supervisor/celery_worker_A_stderr.log
sudo tail -n 50 /var/log/supervisor/celery_worker_log_stderr.log


curl http://127.0.0.1:8083/docs # 或者服务器的公网/内网 IP